# cookbooks/consul_wrapper/templates/default/ruby-configure-consul.conf.erb
respawn

start on runlevel [2345]
stop on runlevel [06]

script
  env <%= @env.map { |e| e.join('=') }.join(' ') %> ruby <%= @script_path %>
end script
